# src/init.d/nftlist
#!/sbin/openrc-run

extra_started_commands="refresh panic"

NFTLIST='nftlist-little.sh'


depend() {
	need nftables
	before ulogd
	after networking
}

start() {
	ebegin "Starting ${RC_SVCNAME}"
	$NFTLIST update
	eend $?
}

stop() {
	ebegin "Stopping ${RC_SVCNAME}"
	$NFTLIST purge
	eend $?
}

refresh() {
	ebegin "Refreshing ${RC_SVCNAME}"
	$NFTLIST update
	eend $?
}

panic() {
	ebegin "PANIC CALL: ${RC_SVCNAME}"
	$NFTLIST panic
	eend $?
}
