# src/nft-helper
#!/sbin/openrc-run

extra_started_commands="refresh"
RULE_DIR="/etc/nftdef"

depend() {
	need nftables
	before ulogd
	after networking
}

load_config() {
	for file in $(ls $RULE_DIR/*.list); do
		TABLE=""
		CHAIN=""
		RULES=""

		# split file name into parts that correspond to:
		# $TABLE-$CHAIN-$RULES.list
		for a in $(echo $file | tr "-" " "); do
			if [ -z $TABLE ] ; then
				# use basename to cutout 'path' string
				TABLE=$(basename $a)
			elif [ -z $CHAIN ] ; then
				CHAIN=$a
			else
				# cut file extension (and dot)
				RULES=$(echo $a | tr "." "\n" | head -1)
			fi
		done

		# init/update/discard the set elements
		nft-helper.sh $1 $TABLE $CHAIN $RULES $file
	done
}

start() {
	load_config 'init'
}

refresh() {
	load_config 'update'
}
